# laravel/8.0/files/etc/services.d/queue/run
#!/usr/bin/with-contenv sh

exec php8 /var/www/html/artisan queue:work
