# laravel/8.0/files/etc/services.d/queue/run
#!/usr/bin/with-contenv sh

exec php /app/artisan queue:work
