# laravel/8.0/files/crontab
* * * * * cd /var/www/html && php8 artisan schedule:run >> /dev/null 2>&1
